# recipes-framework/ollama/files/ollama.path
[Path] 
PathExists=/dev/nvidia-caps

[Install] 
WantedBy=multi.user.target
